# docker/services/cron/cron.d/circulation
# /etc/cron.d/circulation: Palace Manager crontab
# You don't have to run the `crontab' command to install the new
# version when you edit this file in /etc/cron.d. Files in this directory
# also have username fields, similar to the systemwide /etc/crontab.

SHELL=/bin/sh
PATH=/usr/local/sbin:/usr/local/bin:/sbin:/bin:/usr/sbin:/usr/bin
HOME=/var/www/circulation

# m h dom mon dow user command

# These scripts update internal caches.
#
0 0 * * * root bin/run update_custom_list_size >> /var/log/cron.log 2>&1
0 2 * * * root bin/run update_lane_size >> /var/log/cron.log 2>&1
*/30 * * * * root bin/run -d 5 equivalent_identifiers_refresh >> /var/log/cron.log 2>&1

# The remaining scripts keep the circulation manager in sync with
# specific types of collections.

# Bibliotheca
#
*/15 * * * * root bin/run bibliotheca_monitor >> /var/log/cron.log 2>&1
0 */5 * * * root bin/run bibliotheca_purchase_monitor >> /var/log/cron.log 2>&1
0 5 * * * root bin/run bibliotheca_circulation_sweep >> /var/log/cron.log 2>&1

# Overdrive
#
*/15 * * * * root bin/run overdrive_reaper >> /var/log/cron.log 2>&1
0 4 * * * root bin/run overdrive_format_sweep >> /var/log/cron.log 2>&1

# Auto update lists
# Every hour between 7AM - 1AM
5 0,1,7-23 * * * root bin/run custom_list_update_new_entries >> /var/log/cron.log 2>&1
