# docker/runit/uwsgi/run
#!/bin/bash
set -e

exec /var/www/circulation/env/bin/uwsgi --ini /etc/uwsgi.ini
